# docker/crontab
34 * * * * /Preco/src/main/script/health-report-watchdog.sh

